# install/wfs.clean.service
[Unit]
Description=WFS-CLEAN

[Service]
User=pi
Type=simple
ExecStart=/usr/bin/python3 /home/pi/wfs/wfs_sub_db_cleanup.py

